# zfs-linux/fs-sentinel-daemon.service
[Service]
Type=simple
ExecStart=/usr/local/bin/fs-sentinel-zfs-daemon.sh

[Unit]
Description=Daemon for monitoring all mounted ZFS datasets using fs-sentinel
After=zfs-mount.service
Wants=zfs-mount.service

[Install]
WantedBy=default.target
